# supervisord.conf
[supervisord]
logfile = /tmp/supervisord.log
loglevel = DEBUG
nodaemon=true

[program:web]
command=python -m lndngigs.web

[program:worker]
command=python -m lndngigs.worker
